# lastlogin.crontab
0  2  *  *  * root sleep $(($RANDOM\%30))m /etc/zabbix/bin/lastlogin.sh
